# eve/workers/build/supervisord.conf
[program:start_dmd]
command=npm run start_dmd
directory=/home/eve/S3
user=eve
autostart=true
autorestart=false
process_name=dmd
numprocs=1
stdout_logfile=/tmp/dmd_server.txt
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/dmd_server.txt
stderr_logfile_maxbytes=0
stopsignal=QUIT
#environment=
#    FOO=bar
#    FOO2=bar
