# eve/workers/build/supervisord.conf
#[program:start_dmd]
#command=/bin/bash -l /home/eve/S3/tests.bash
#directory=/home/eve/S3
#autostart=true
#autorestart=true
#process_name=dmd
#numprocs=1
#stdout_logfile=/dev/fd/1
##stdout_logfile=/tmp/dmd_server.txt
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/fd/2
##stderr_logfile=/tmp/dmd_server.txt
#stderr_logfile_maxbytes=0
#stopsignal=QUIT
##environment=
##    FOO=bar
##    FOO2=bar

[program:redis]
command=/usr/bin/redis-server
autostart=true
autorestart=true
